# deploy/systemd/opaire-disk-cleanup.timer
[Unit]
Description=Run opaiRe disk cleanup periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
